# nas/templates/default/system/rc.local.erb
#!/bin/sh
# Generated by Chef for <%= node['fqdn'] %>

btrfs device scan || true
mount /storage

exit 0
