# nas/templates/default/system/rc.local.erb
#!/bin/sh
# Generated by Chef for <%= node['fqdn'] %>

<% Dir.glob("/dev/sd[a-z]").each do |drive| -%>
hdparm -S 244 <%= drive %>
<% end -%>

exit 0
