# nas/templates/default/system/rc.local.erb
#!/bin/sh
# Generated by Chef for <%= node['fqdn'] %>

<% Dir.glob("/dev/sd[a-z]").each do |drive| -%>
hdparm -S 244 <%= drive %>
<% end -%>

if [ -e /dev/disk/by-label/backup -a -e /backup ]; then
    mount /backup
fi

exit 0
