# nas/templates/default/system/rc.local.erb
#!/bin/sh -e
# Generated by Chef for <%= node['fqdn'] %>

btrfs device scan
mount /storage

exit 0
