# nas/templates/default/system/rc.local.erb
#!/bin/bash
# Generated by Chef for <%= node['fqdn'] %>

<% Dir.glob("/dev/sd[a-z]").each do |drive| -%>
hdparm -S 244 <%= drive %>
<% end -%>

if modprobe -s bcache
then
    for i in {a..e} 
    do 
	echo /dev/sd$i > /sys/fs/bcache/register_quiet
    done
    sleep 3
    echo 64M > /sys/block/bcache0/bcache/sequential_cutoff
fi

if [ -e /dev/disk/by-label/storage -a -e /storage ]; then
    mount /storage
    sudo chmod 775 /storage/
    sudo chown kimata:www-data /storage/
fi
if [ -e /dev/disk/by-label/backup -a -e /backup ]; then
    mount /backup
fi

exit 0
