# router/templates/default/system/rc.local.erb
../../../../nas/templates/default/system/rc.local.erb